# crontab
0 3 * * * /usr/local/bin/log-archive /var/log >> /var/log/cron.log 2>&1
